# deploy/babytrack.service
[Unit]
Description=BabyTrack API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=babytrack
Group=babytrack
WorkingDirectory=/opt/babytrack
ExecStart=/opt/babytrack/babytrack -config /opt/babytrack/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=babytrack

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/babytrack

# Environment
Environment=GIN_MODE=release

[Install]
WantedBy=multi-user.target
