# deploy/family-tracker.service
[Unit]
Description=Family Tracker API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=family-tracker
Group=family-tracker
WorkingDirectory=/opt/family-tracker
ExecStart=/opt/family-tracker/family-tracker -config /opt/family-tracker/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=family-tracker

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/family-tracker

# Environment
Environment=GIN_MODE=release

[Install]
WantedBy=multi-user.target
